# cccbot-scraper.service
[unit]
Description=Run CCC bot scraper
After=network.target

[Service]
User=webmaster
WorkingDirectory=/home/webmaster/ccc-bot/
ExecStart=/home/webmaster/ccc-bot/start-scraper.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
